# plugin/app2sd/app2sd-server.service.in
[Unit]
Description=app2sd server

[Service]
ExecStart=@PREFIX@/bin/app2sd-server
CapabilityBoundingSet=~CAP_MAC_ADMIN
CapabilityBoundingSet=~CAP_MAC_OVERRIDE
